# device/main/device-purism-librem5/70-l5-xhci-bus-recover.rules
SUBSYSTEM=="platform", DEVPATH=="/devices/platform/soc@0/38200000.usb/xhci-hcd.4.auto", ACTION=="remove", RUN+="/usr/sbin/l5-xhci-reset"
