# main/ts-srv/ts-srv.init
#!/sbin/openrc-run

command="/usr/sbin/ts_srv"
supervisor="supervise-daemon"

depend() {
	after bootmisc
}
